# testing/heartbeat/eb_heartbeat.timer
[Unit]
Description=Record heartbeat every 10s

[Timer]
OnBootSec=1min
OnUnitActiveSec=10s

[Install]
WantedBy=timers.target