# testing/heartbeat/eb_heartbeat.timer
[Unit]
Description=Record heartbeat every 10s

[Timer]
OnBootSec=10s
OnUnitActiveSec=10s
Persistent=true

[Install]
WantedBy=timers.target